# webclaw.service
[Unit]
Description=WebClaw - AI Agent Sessions in Firecracker microVMs
After=network.target

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/webclaw
ExecStart=/home/exedev/webclaw/webclaw
Restart=on-failure
RestartSec=5
Environment=HOME=/home/exedev

[Install]
WantedBy=multi-user.target
